# etc/systemd/system/mnt-data.automount
[Unit]
Description=Media on Toshilab

[Automount]
Where=/mnt/data

[Install]
WantedBy=multi-user.target
